# config/gitea-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gitea.act-runner</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/Vinod/Documents/GitHub/community-rule/act_runner</string>
        <string>daemon</string>
        <string>--config</string>
        <string>/Users/Vinod/Documents/GitHub/community-rule/config/gitea-runner.yaml</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/Vinod/Documents/GitHub/community-rule</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/Vinod/Documents/GitHub/community-rule/runner.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/Vinod/Documents/GitHub/community-rule/runner-error.log</string>
</dict>
</plist>
